# Docker/docker-compose/xray-agent/xray/backup_crontab.cron
30 1 * * * /bin/bash /etc/v2ray-agent/install.sh RenewTLS